# modules/example/templates/celerybeat-supervisor.conf.erb
[program:celerybeat]
directory=/opt/celery
command=celerybeat --loglevel DEBUG
stdout_logfile=/var/log/celerybeat.log
stderr_logfile=/var/log/celerybeat.log
