# modules/example/templates/celeryd-supervisor.conf.erb
[program:celerybeat]
directory=/opt/celery
command=celeryd --loglevel DEBUG --events
stdout_logfile=/var/log/celeryd.log
stderr_logfile=/var/log/celeryd.log
